# braindump.service
[Unit]
Description=Braindump Application
After=network.target postgresql.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=alex
WorkingDirectory=/home/alex/braindump
ExecStartPre=/usr/bin/npm install --prefix backend
ExecStart=/usr/bin/node backend/server.js
Restart=always
RestartSec=10

# Exponential backoff
RestartSteps=10
RestartMaxDelaySec=300

# Environment
Environment=NODE_ENV=production
Environment=PORT=8001

# Resource limits
MemoryLimit=512M
CPUQuota=80%

# Health monitoring
ExecStartPost=/bin/sleep 5
ExecStartPost=/usr/bin/curl -f http://localhost:8001/health || exit 1

# Watchdog (requires systemd 229+)
WatchdogSec=30
NotifyAccess=all

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=braindump

[Install]
WantedBy=multi-user.target